# claude-triager/lxc/git-pull.service
[Unit]
Description=Pull latest homeops repo for claude-triager (read-only)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=triager
Group=triager
WorkingDirectory=/opt/claude-triager/homeops
Environment=HOME=/var/lib/claude-triager
Environment=GIT_SSH_COMMAND=ssh -i /var/lib/claude-triager/.ssh/id_ed25519 -o IdentitiesOnly=yes
ExecStart=/usr/bin/git -C /opt/claude-triager/homeops pull --ff-only --quiet
# Never allow a push from this unit
ExecStartPre=/usr/bin/git -C /opt/claude-triager/homeops config --local receive.denyCurrentBranch true
