# claude-triager/lxc/triager.service
[Unit]
Description=Claude SRE Triager - Autonomous Kubernetes alert investigator
Documentation=man:systemd.service(5)
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
# Basic service configuration
Type=simple
User=triager
Group=triager
WorkingDirectory=/opt/claude-triager

# Environment and execution
EnvironmentFile=/etc/claude-triager.env
Environment="PATH=/opt/claude-triager/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/opt/claude-triager/venv/bin/python3 /opt/claude-triager/bot/main.py

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/var/log/triager.log
StandardError=append:/var/log/triager.log
SyslogIdentifier=triager

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictNamespaces=yes
RestrictSUIDSGID=yes
RemoveIPC=yes
PrivateDevices=yes
ProtectClock=yes
ProtectHostname=yes
LockPersonality=yes
SystemCallFilter=@system-service

# Paths that must be writable
ReadWritePaths=/opt/claude-triager/state /var/log/triager.log /var/lib/claude-triager

[Install]
WantedBy=multi-user.target