# kubernetes/apps/monitoring/exporters/wireguard-oci/host/wg-exporter.service
[Unit]
Description=WireGuard Prometheus exporter (read-only, serves :9586)
After=network-online.target wg-quick@oci.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/wg-exporter.py
Restart=on-failure
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
